# conf/etc/service/carbon/run
#!/bin/sh

exec python3 /opt/graphite/bin/carbon-cache.py start $CARBON_OPTS 2>&1
