# conf/etc/service/carbon-aggregator/run
#!/bin/sh

exec python3 /opt/graphite/bin/carbon-aggregator.py start $CARBON_OPTS 2>&1
