# conf/etc/service/carbon-relay/run
#!/bin/sh

[[ -n "${RELAY}" ]] || exit 1

exec python3 /opt/graphite/bin/carbon-relay.py start $CARBON_OPTS 2>&1
